# systemd/battery-health.service
[Unit]
Description = Battery Health Monitoring Service

[Service]
Type = oneshot
ExecStart = /usr/bin/python3 /opt/battery_health/app.py
