# systemd/battery-health.timer
[Unit]
Description = Run Battery Health Monitor every 5 minutes

[Timer]
OnBootSec = 2min
OnUniActiveSec = 5min
AccuracySec = 30s

[Install]
WantedBy = timers.target
